# processes_and_signals/11-manage_my_process
#!/usr/bin/env bash
#!/bin/bash
### BEGIN INIT INFO
# Provides:          manage_my_process
# Required-Start:    $local_fs $network $syslog
# Required-Stop:     $local_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Manages the manage_my_process script
### END INIT INFO
# Change the path to your manage_my_process script
manage_my_process_script="https://github.com/Samenergy/alu-shell/blob/14ecee31c2197f78f91bcf3e85b651b6f65454e2/processes_and_signals/manage_my_process"
pid_file="/var/run/my_process.pid"
case "$1" in
    start)
        if [ -f "$pid_file" ]; then
            echo "manage_my_process is already started"
            exit 1
        fi
        /bin/bash "$manage_my_process_script" > /tmp/my_process &
        echo $! > "$pid_file"
        echo "manage_my_process started"
        ;;
    stop)
        if [ ! -f "$pid_file" ]; then
            echo "manage_my_process is not running"
            exit 1
        fi
        kill "$(cat "$pid_file")"
        rm "$pid_file"
        echo "manage_my_process stopped"
        ;;
    restart)
        if [ -f "$pid_file" ]; then
            kill "$(cat "$pid_file")"
            rm "$pid_file"
        fi
        /bin/bash "$manage_my_process_script" > /tmp/my_process &
        echo $! > "$pid_file"
        echo "manage_my_process restarted"
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
esac

exit 0

